# dockerconf/crontab.cron
* * * * * php /data/wwwroot/kael/yii ldap/update